# scripts/crontab
@hourly /usr/local/bin/musicbox_kafka_collector.py
@daily /usr/local/bin/musicbox_aggregate_daily_log.sh

